# kiosk/kittenbus.cron
# Cron file for Kitten bus

# Crontab fields:
#
# field #   meaning        allowed values
# -------   ------------   --------------
#    1      minute         0-59
#    2      hour           0-23
#    3      day of month   1-31
#    4      month          1-12 (or names, see below)
#    5      day of week    0-7 (0 or 7 is Sun, or use names)

# Restart every night at 4AM.  The shutdown command needs sudo
0 4 * * * osascript -e 'tell app "System Events" to restart'
